# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

; ── Docker daemon ──────────────────────────────────────────────────────────────
[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock --host=tcp://0.0.0.0:2375
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/dockerd.log
stderr_logfile=/var/log/dockerd.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Compose launcher ───────────────────────────────────────────────────────────
[program:compose-up]
command=/start.sh
autostart=true
autorestart=false
priority=10
stdout_logfile=/var/log/compose-up.log
stderr_logfile=/var/log/compose-up.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ── Keep-alive (holds the container open after compose detaches) ───────────────
[program:keepalive]
command=tail -f /dev/null
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/null
stderr_logfile=/dev/null
